# overlay/2020/etc/update-motd.d/99-mev
#!/bin/sh

findmnt -t ext4
findmnt -t overlay
ls -al /sbin/init

